# pi/assistant.service
[Unit]
Description=ESP32 Assistant Brains (Audio & Rotary)
After=network.target

[Service]
# Mute the amp immediately — before Python starts — so the SD pin is never
# floating while the I2S clock and ALSA drivers initialise
ExecStartPre=/usr/bin/pinctrl set 13 op dl
# Use the Python from our virtual environment
ExecStart={{APP_PATH}}/venv/bin/python -u {{APP_PATH}}/assistant_brains.py
WorkingDirectory={{APP_PATH}}
# Load environment variables (API Keys, etc.) from .env
EnvironmentFile=-{{APP_PATH}}/.env
StandardOutput=inherit
StandardError=inherit
Restart=always
User={{SERVICE_USER}}
Group={{SERVICE_USER}}

[Install]
WantedBy=multi-user.target
